# systemd/deploy.service
[Unit]
Description=GitHub Auto-Deploy Webhook Server
After=network.target

[Service]
Type=simple
# User=YOUR_USER
# Group=YOUR_GROUP
WorkingDirectory=/path/to/deploy
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
# EnvironmentFile=/path/to/deploy/.env

[Install]
WantedBy=multi-user.target
